# magento2/manager/phusion/crontab
*/1 * * * * root /usr/bin/php /home/dev/PAJ/www/Docker/Scale/Manager/cron.php silent scalemanager

